# systemd/.config/systemd/user/foot-server.service
[Unit]
Description=foot server

[Service]
ExecStart=foot --server
Restart=on-failure

[Install]
WantedBy=default.target

